# deploy/systemd/screenfish-docker.user.service
[Unit]
Description=ScreenFish (Docker Compose)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/services/screenfish-prod/backend

# Use explicit env files:
# - `.env.docker` for ports and non-secret settings
# - (optional) another env file for secrets (auth/smtp/admin token)
ExecStartPre=/bin/bash -lc 'for i in {1..30}; do /usr/bin/docker info >/dev/null 2>&1 && exit 0; sleep 1; done; echo \"docker not ready\" >&2; exit 1'
ExecStart=/bin/bash -lc 'set -euo pipefail; BASE=%h/services/screenfish-prod; export SCREENFISH_WEB_CONTEXT=../frontend; export SCREENFISH_DATA_DIR=%h/workspace/screenfish/data; export STOCK_SCREENER_GIT_SHA="$(git --git-dir=${BASE}/git/screenfish.git rev-parse --short refs/heads/main 2>/dev/null || true)"; export STOCK_SCREENER_GIT_DESCRIBE="$(git --git-dir=${BASE}/git/screenfish.git describe --tags --always refs/heads/main 2>/dev/null || true)"; export STOCK_SCREENER_BUILD_TIME="$(date -u +%Y-%m-%dT%H:%M:%SZ)"; export VITE_BUILD_SHA="$(git --git-dir=${BASE}/git/screenfish-web.git rev-parse --short refs/heads/main 2>/dev/null || true)"; export VITE_APP_VERSION="$(git --git-dir=${BASE}/git/screenfish-web.git describe --tags --always refs/heads/main 2>/dev/null || true)"; export VITE_BUILD_TIME="${STOCK_SCREENER_BUILD_TIME}"; /usr/bin/docker compose --env-file %h/workspace/screenfish/.env.docker --env-file ${BASE}/run/backend.env up -d --remove-orphans --build'
ExecStop=/bin/bash -lc 'set -euo pipefail; export SCREENFISH_WEB_CONTEXT=../frontend; export SCREENFISH_DATA_DIR=%h/workspace/screenfish/data; /usr/bin/docker compose --env-file %h/workspace/screenfish/.env.docker --env-file %h/services/screenfish-prod/run/backend.env down'

RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=default.target
